# rds_nginx/supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info


[program:nginx]
command=nginx -g 'daemon off;'

[program:fcgiwrap]
command=fcgiwrap -s unix:/var/run/fcgiwrap.socket